# hutwatch.service
[Unit]
Description=HutWatch BLE Temperature Monitor
After=network.target bluetooth.target local-fs.target
Wants=bluetooth.target

[Service]
Type=simple
User=YOURUSER
Group=YOURUSER
WorkingDirectory=/home/YOURUSER/hutwatch

# Capabilities for Bluetooth adapter reset without sudo
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

ExecStart=/home/YOURUSER/hutwatch/venv/bin/python -m hutwatch -c config.yaml

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
