# hutwatch.service
[Unit]
Description=HutWatch BLE Temperature Monitor
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=tronkko
Group=tronkko
WorkingDirectory=/home/tronkko/code/hutwatch
ExecStart=/home/tronkko/code/hutwatch/venv/bin/python -m hutwatch -c config.yaml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
